# etc/init.d/vzeventd-gentoo.in
#!/sbin/runscript

depend() {
	use net
	before vz
}

start() {
	ebegin "Loading vzevent module"
	modprobe -r vzevent >/dev/null 2>&1
	modprobe vzevent reboot_event=1 >/dev/null 2>&1
	rc=$?
	if [ ${rc} -ne 0 ]; then
		eend ${rc} "failed"
		return ${rc}
	fi
	eend ${rc}

	ebegin "Starting vzeventd"
	start-stop-daemon --start --exec /usr/sbin/vzeventd 
	eend $?
}

stop() {
	ebegin "Stopping vzeventd"
	start-stop-daemon --stop -n vzeventd 
	eend $?
}
